# root/etc/s6-overlay/s6-rc.d/unbound-anchor/up
s6-setuidgid unbound
unbound-anchor -a /etc/unbound/root.key
